# glusterd/overlay/etc/service/glustereventsd/run
#!/usr/bin/env bash

export PYTHONPATH=/usr/lib/python2.7/dist-packages:$PYTHONPATH

exec /usr/sbin/glustereventsd --pid-file /var/run/glustereventsd.pid